# scraper/crontab
* * * * 1-5 python /task.py >> /logfile
  